# systemd/nsaf-orchestrator.service
[Unit]
Description=NSAF Orchestrator — Nightshift AutoFoundry queue and session manager
After=network.target postgresql.service

[Service]
Type=simple
User=NSAF_USER
Group=NSAF_USER
WorkingDirectory=NSAF_DIR
ExecStart=/usr/bin/node orchestrator/src/index.js
Restart=on-failure
RestartSec=10
EnvironmentFile=NSAF_DIR/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nsaf-orchestrator

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
