# systemd/nsaf-flask.service
[Unit]
Description=NSAF Flask App — Webex bot, idea selection UI, QA review
After=network.target

[Service]
Type=simple
User=NSAF_USER
Group=NSAF_USER
WorkingDirectory=NSAF_DIR
ExecStart=NSAF_DIR/venv/bin/python flask-app/app.py
Restart=on-failure
RestartSec=10
EnvironmentFile=NSAF_DIR/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nsaf-flask

[Install]
WantedBy=multi-user.target
